# images/fedora43-asus-tuf-a16/files/dreamcatcher-oem-firstboot.service
[Unit]
Description=Dreamcatcher OEM first-boot tasks
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/var/lib/dreamcatcher-oem/firstboot-complete

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/dreamcatcher-oem-firstboot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
